# evo.service
# systemd service — auto-start Evo on boot
# Install:
#   sudo cp evo.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now evo
#
# Logs:  journalctl -u evo -f

[Unit]
Description=Evo - Self-Evolving Elixir Application
After=network.target

[Service]
Type=simple
User=tuffypi
Group=tuffypi
WorkingDirectory=/home/tuffypi/Projects/claude_brain
EnvironmentFile=/home/tuffypi/Projects/claude_brain/.env
Environment=MIX_ENV=prod
ExecStart=/bin/bash -lc 'cd /home/tuffypi/Projects/claude_brain && exec elixir --sname evo -S mix phx.server'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
